# dockers/grafana/etc/supervisor.d/nginx.ini
[program:nginx]
autostart = true
autorestart = true
stdout_logfile=/grafana/logs/grafana-nginx.log
stderr_logfile=/grafana/logs/grafana-nginx-error.log
stdout_events_enabled = false
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0

command = /usr/sbin/nginx -c /etc/nginx/nginx.conf
